# steve/rootfs/etc/s6-overlay/s6-rc.d/start-steve/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant App: SteVe
# Runs the SteVe OCPP Server
# ==============================================================================

bashio::log.info "Starting SteVe...."
exec /usr/lib/jvm/default-jvm/jre/bin/java -jar /data/target/steve.war